# root/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

# Log level 5 (and below) is noisy during periodic wakeup where nothing happens.
/usr/sbin/crond -f -l 6
